# tools/cron/task.cron
# Push changes on ~/.task
*/20 * * * * cd ~/.task && git add * && git commit -am update && git push origin 2>&1
