# tools/cron/task.cron
# Push changes on ~/.task
0 * * * * (cd ~/.task && git add * && git commit -am update && git push origin &>/dev/null &)
