# script/sync-vimwiki.cron
0 * * * * cd ~/vimwiki && git commit -am "update" && git push
