# scripts/sync-vimwiki.cron
0 * * * * cd ~/vimwiki && git add * && git commit -am "update" && git push
